# dockerfiles/etc/service/nginx/run
#!/bin/sh
set -e
exec /usr/local/openresty/nginx/sbin/nginx -g "daemon off;"